# deploy/production/systemd/worldcup-refresh-daily-00.timer
[Unit]
Description=Run World Cup daily 00:00 refresh

[Timer]
OnCalendar=*-*-* 00:05:00
Persistent=true
Unit=worldcup-refresh@daily_00.service

[Install]
WantedBy=timers.target
